# debian/postinst

#!/bin/sh
set -e
pip3 install paho-mqtt psutil openpyxl
systemctl enable monitor.service
systemctl start monitor.service
exit 0
