# debian/prerm
#!/bin/sh
set -e

systemctl stop mqtt-monitor.service
systemctl disable mqtt-monitor.service

exit 0
